# ansible_files/roles/backend/files/node_exporter/node-exporter.service
[Unit]
Description=node-exporter service

[Service]
Type=simple
Restart=on-failure
ExecStart=/node_exporter/node_exporter

[Install]
WantedBy=multi-user.target

